# salt/roots/salt/demo-project/celery.conf
{% set workers = pillar.get('celery', {}).get('workers', []) %}
{%- for worker in workers %}
[program:{{ worker['name'] }}]
command={{ pillar['django']['virtualenv'] }}/bin/python
		{{ pillar['django']['virtualenv'] }}/bin/celery
		worker
		--app {{ pillar['django']['project'] }}
		--workdir={{ pillar['django']['path'] }}
		--pidfile=/var/run/demo-project/{{ worker['name'] }}.pid
	{%- for option,value in worker.options.items() %}
	  {%- if not value %}
		--{{ option }}
		{%- else %}
		--{{ option }}={{ value }}
		{%- endif %}
  {%- endfor %}
user={{ pillar['django']['user'] }}
directory={{ pillar['django']['path'] }}
environment=PYTHONPATH={{ pillar['django']['path'] }}
killasgroup=true
priority=998
startsecs=10
stderr_logfile={{ pillar['django']['log_path'] }}/{{ worker['name'] }}.log
stdout_logfile={{ pillar['django']['log_path'] }}/{{ worker['name'] }}.log
stopwaitsecs=600
{%- endfor %}